# task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Reply Guy Monitor - hourly tweet scraper, runs during your active hours</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T08:00:00</StartBoundary>
      <Repetition>
        <Interval>PT1H</Interval>
        <Duration>PT17H</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT10M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <!-- Update these three paths before importing the task -->
      <Command>C:\path\to\python.exe</Command>
      <Arguments>C:\path\to\reply-guy-monitor\scripts\x_feed_monitor.py</Arguments>
      <WorkingDirectory>C:\path\to\reply-guy-monitor</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
